# office/config/systemd/user/dawarich-backup.timer
[Unit]
Description=Dawarich backup timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=default.target
