# office/config/systemd/user/paperless-backup.timer
[Unit]
Description=Paperless backup timer

[Timer]
OnCalendar=weekly
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=default.target
