# assets/podman/paperless/postgres.container
# vim: ft=systemd

[Unit]
Description=PostgreSQL

[Container]
Pod=paperless.pod
ContainerName=postgres
Image=docker.io/library/postgres:18
UserNS=keep-id:uid=1000,gid=1000
Environment=POSTGRES_INITDB_ARGS=--data-checksums
Volume=/data/paperless/postgres:/var/lib/postgresql:Z
HealthCmd=pg_isready -U postgres
HealthStartPeriod=5m
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5
Notify=healthy
Label="io.containers.autoupdate=registry"

[Service]
Restart=always
