# assets/systemd/backup-failure-notify@.service
[Unit]
Description=Notify on %i failure

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'journalctl --user -u %i --no-pager -n 50 | /opt/scripts/send-notification.sh "[%H] %i failed"'
