# assets/systemd/immich-backup.timer
[Unit]
Description=Immich backup timer

[Timer]
OnCalendar=weekly
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
