# assets/fs/etc/systemd/system/immich-backup.service
[Unit]
Description=Immich backup job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/backup/immich.sh
