# office/config/containers/systemd/minuspod/minuspod.container
# vim: ft=systemd

[Unit]
Description=MinusPod podcast ad removal

[Container]
ContainerName=minuspod
Image=docker.io/ttlequals0/minuspod:latest
AddDevice=nvidia.com/gpu=all
Network=minuspod.network
Environment=BASE_URL=https://podcasts.guirao.net
Environment=SESSION_COOKIE_SECURE=false
Environment=WHISPER_MODEL=large-v3-turbo
Environment=WHISPER_COMPUTE_TYPE=int8_float16
Environment=WHISPER_DEVICE=cuda
Environment=LLM_PROVIDER=openai-compatible
Environment=OPENAI_MODEL=ad-detect
Environment=OPENAI_BASE_URL=http://litellm:4000/v1
EnvironmentFile=%E/minuspod/env
Volume=%S/minuspod/data:/app/data:Z
Volume=/etc/localtime:/etc/localtime:ro
UserNS=keep-id:uid=1000,gid=1000
PublishPort=10.10.10.11:4444:8000
HealthCmd=curl -f http://localhost:8000/api/v1/health
HealthInterval=30s
HealthTimeout=3s
HealthStartPeriod=20s
HealthRetries=3
Label="io.containers.autoupdate=registry"

[Service]
Restart=always

[Install]
WantedBy=default.target
